# deploy/systemd/sql2doc.service
# Systemd service file for SQL2Doc
# Copy to: /etc/systemd/system/sql2doc.service

[Unit]
Description=SQL2Doc Application (Docker Compose)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/sql2doc
User=azureuser
Group=docker

# Load environment variables
EnvironmentFile=/opt/sql2doc/.env

# Start services
ExecStart=/usr/bin/docker-compose up -d

# Stop services
ExecStop=/usr/bin/docker-compose down

# Reload services
ExecReload=/usr/bin/docker-compose restart

# Restart policy
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target